# crontab
* * * * * root netstat  -pln | grep 5672 || /sbin/service rabbitmq-server restart
